# scripts/orion_kalshi_autotrade_cycle_launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.orion.kalshi_autotrade_cycle</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/kalshi_autotrade_cycle_run.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>300</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/corystoner/Library/Logs/orion_kalshi_autotrade_cycle.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/corystoner/Library/Logs/orion_kalshi_autotrade_cycle.log</string>
</dict>
</plist>
