# scripts/orion_resurrector_launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.openclaw.orion.resurrector</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>-lc</string>
      <string>"/bin/bash -lc '__REPO_ROOT__/scripts/resurrect_orion_mac.sh'"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>${HOME}/Library/Logs/orion_resurrector.log</string>
    <key>StandardErrorPath</key>
    <string>${HOME}/Library/Logs/orion_resurrector.log</string>
  </dict>
</plist>
